# tuya-update.timer
[Unit]
Description=Update Tuya thermometer data every minute

[Timer]
OnBootSec=10s
OnUnitActiveSec=120s

[Install]
WantedBy=timers.target
